# ndi-streamer.service
[Unit]
Description=NDI Streamer Service
Documentation=https://github.com/YOUR_ORG/ndi-streamer
After=network.target
Wants=network.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=ndi-streamer
Group=ndi-streamer
WorkingDirectory=/opt/ndi-streamer
Environment=PATH=/opt/ndi-streamer/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=LD_LIBRARY_PATH=/usr/share/ndi/lib
# ndi-streamer runs with --no-create-home and ProtectHome=true, so $HOME is
# unwritable. Pin Playwright's browser cache to a path the service can reach.
Environment=PLAYWRIGHT_BROWSERS_PATH=/opt/ndi-streamer/ms-playwright
Environment=HOME=/opt/ndi-streamer
EnvironmentFile=/opt/ndi-streamer/.env

ExecStart=/opt/ndi-streamer/venv/bin/python run.py

Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ndi-streamer

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/ndi-streamer
ProtectHome=true
PrivateTmp=true
UMask=0077
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
# AF_NETLINK is required for mDNS / NDI discovery on Linux.
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK

# Resource limits
LimitNOFILE=65536
MemoryMax=4G
TasksMax=4096

[Install]
WantedBy=multi-user.target
